# bookstack-rag/rootfs/etc/s6-overlay/s6-rc.d/api/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Halt the add-on if the API service exits, so HA Supervisor can restart it.
# ==============================================================================

bashio::log.warning "API service exited — bringing the add-on down."

exec /run/s6/basedir/bin/halt
