# bookstack-rag/rootfs/etc/s6-overlay/s6-rc.d/api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the BookStack RAG FastAPI server, after waiting for the qdrant
# sidecar to accept connections on localhost:6333.
# ==============================================================================

bashio::log.info "Waiting for qdrant readiness on localhost:6333..."
for attempt in $(seq 1 60); do
    if curl -sfm 1 http://localhost:6333/readyz >/dev/null 2>&1; then
        bashio::log.info "qdrant ready after ${attempt}s — starting API"
        break
    fi
    sleep 1
done

cd /app

exec python3 -m uvicorn \
    --factory app.main:create_app \
    --host 0.0.0.0 \
    --port 8000
